# hello-world-backend.service
[Unit]
Description=Hello World Backend Service
Documentation=https://github.com/yourusername/hello-world
After=network.target mysql.service

[Service]
Type=simple
User=www
Group=www
WorkingDirectory=/www/wwwroot/hello-world/backend
Environment="PATH=/usr/bin:/usr/local/bin"
EnvironmentFile=/www/wwwroot/hello-world/backend/.env

# 使用gunicorn运行Flask应用
ExecStart=/usr/bin/python3 -m gunicorn \
    --workers 2 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile /var/log/hello-world-access.log \
    --error-logfile /var/log/hello-world-error.log \
    app:app

# 自动重启配置
Restart=always
RestartSec=5
StartLimitInterval=0

# 安全配置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
